# worker/setup-tunnel.service
[Unit]
After=network.target

[Service]
ExecStart=/root/gerrit-ci-scripts/worker/tunnel.sh

[Install]
WantedBy=default.target
